# rag-service.service.example
[Unit]
Description=RAG Microservice - AI-powered document Q&A
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/rag-service
Environment="PATH=/var/www/rag-service/venv/bin"
ExecStart=/var/www/rag-service/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=append:/var/www/rag-service/logs/systemd.log
StandardError=append:/var/www/rag-service/logs/systemd-error.log

[Install]
WantedBy=multi-user.target
